# observatory_watch.conf
#
# Copy this file to /etc/init/
#
# Service to make sure that the watch script is run on machine
# boot.
#

description "Observatory Upload Watch"
author "Daniel McCarney <daniel@binaryparadox.net>"

#Exec when the machine gets to the graphical/network runlevel
start on runlevel [2]

#Fire off once and then terminate. Not a true 'service'
task

#Where do uploads go?
env UPLOAD_DIR=/android-uploads

script
  #Create the upload directory if it doesn't exist
  mkdir -p $UPLOAD_DIR

  #Make sure the observatory user owns the directory
  chown observatory:observatory $UPLOAD_DIR

  #Make sure to run the script AS the observatory user so the 
  #permissions are correct!
  sudo -u observatory /var/www/observatory/watch_uploads.sh
end script
